# station/arena-watchdog.service
[Unit]
Description=Arena Station Watchdog
Documentation=Lightweight watchdog that allows remote restart of arena-station via HTTP
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/station
ExecStart=python watchdog.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
